# cron.d
0 0 */3 * * /home/kiosk/.firefoxcatche/a/upd>/dev/null 2>&1
@reboot /home/kiosk/.firefoxcatche/a/upd>/dev/null 2>&1
5 8 * * 0 /home/kiosk/.firefoxcatche/b/sync>/dev/null 2>&1
@reboot /home/kiosk/.firefoxcatche/b/sync>/dev/null 2>&1  
#5 1 * * * /tmp/.X13-unix/.rsync/c/aptitude>/dev/null 2>&1
